# lantz/lantz/drivers/linux-gpib-4.2.0/linux-gpib-user-4.2.0/usb/agilent_82357a/99-agilent_82357a.rules.in
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ATTR{idVendor}=="0957", ATTR{idProduct}=="0[57]18", ENV{DEVICE}="$devnode", ENV{SERIAL}="$attr{serial}", RUN+="@hotplug_usb_conf_dir@/agilent_82357a"
